# user_scripts/battery/notify/battery_notify.service
[Unit]
Description=Battery Status Notification Service (Event-Driven)
Documentation=man:upower(1)
After=graphical-session.target
PartOf=graphical-session.target

# Require upower to be available
Wants=upower.service
After=upower.service

[Service]
Type=simple
ExecStart=%h/user_scripts/battery/notify/battery_notify.sh

# Restart policy
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=0

# Resource limits (this script is lightweight)
MemoryMax=50M
CPUQuota=5%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=battery-notify

# Security hardening (optional, comment out if issues arise)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
